# launchd/com.nanoclaw.openclaw-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.openclaw-gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/PATH/TO/alfred/scripts/openclaw-gateway.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/PATH/TO/alfred</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/PATH/TO/alfred/logs/openclaw-gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/PATH/TO/alfred/logs/openclaw-gateway.log</string>
</dict>
</plist>
